# docker/supervisor/data/supervisor/conf.d/demo.conf
[program:AudioEncryptor]
command= /data/go/AudioEncryptor/AudioEncryptor --c /data/go/AudioEncryptor/AudioEncryptor.yaml
stopsignal=TERM
autostart = true
autorestart=true
startsecs = 5
user =root
redirect_stderr = true
stdout_logfile = /data/logs/supervisord/AudioEncryptor.log
stdout_logfile_maxbytes=500MB   ; max # logfile bytes b4 rotation (default 50MB)
stdout_logfile_backups=20     ; # of stdout logfile backups (default 10)
stderr_logfile=/data/logs/supervisord/AudioEncryptor.log        ; stderr log path, NONE for none; default AUTO
stderr_logfile_maxbytes=500MB   ; max # logfile bytes b4 rotation (default 50MB)
stderr_logfile_backups=20     ; # of stderr logfile backups (default 10)
